# s6/dovecot/run
#!/bin/sh

exec /usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F