# crontab
0 * * * * root /usr/local/bin/imaputils.pl --host localhost --user catchall --password password --box INBOX --delete --sentbefore $(($(date +%s)-86400))
